# shila/provisioning/conf/etc/init/vagrant-mounted.conf
start on vagrant-mounted
exec sudo service php5-fpm restart
exec sudo service nginx restart
